# mapvideo-uvicorn.service
# Systemd service файл для MapVideo (FastAPI / uvicorn)
# Скопируйте в /etc/systemd/system/mapvideo-uvicorn.service
# Затем:
#   sudo systemctl daemon-reload
#   sudo systemctl enable mapvideo-uvicorn
#   sudo systemctl start mapvideo-uvicorn
#
# При необходимости поправьте пути к conda, каталог проекта и порт (--port 5555).

[Unit]
Description=MapVideo FastAPI (uvicorn)
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/maps
# Предполагается, что conda установлена в /root/miniconda3 и окружение называется "maps"
Environment="PATH=/root/miniconda3/envs/maps/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/miniconda3/envs/maps/bin/uvicorn backend.app.main:app --host 0.0.0.0 --port 5555 --workers 1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=mapvideo.target

